# xyz.openbmc_project.PFR.Manager.service
[Unit]
Description=Daemon manages PFR dbus properties and Event logs

[Service]
Restart=always
ExecStart=/usr/bin/pfr-manager
RestartSec=5
StartLimitInterval=0
SyslogIdentifier=pfr-manager
Type=dbus
BusName=xyz.openbmc_project.PFR.Manager

[Install]
WantedBy=multi-user.target
